# portal/embed/wifi-portal.service
[Unit]
Description=Kazuha Hub WiFi Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 安装时 useradd -r -s /usr/sbin/nologin -d /var/lib/wifi-portal wifi-portal
User=wifi-portal
Group=wifi-portal

# 配置 + 二进制位置. 改了路径同时改这里.
EnvironmentFile=/etc/wifi-portal/wifi-portal.env
ExecStart=/usr/local/bin/wifi-portal

Restart=on-failure
RestartSec=5s

# Portal 已经处理 SIGTERM (graceful shutdown: 停 server + flush banHistory + close EventLog)
KillSignal=SIGTERM
TimeoutStopSec=15s

# /var/lib/wifi-portal 是默认 DATA_DIR, install 时 mkdir + chown wifi-portal
# 跟 .env 里的 DATA_DIR 必须一致.

# === 沙盒 / hardening ===
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/wifi-portal
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
# Portal 不需要 raw socket / device access — 只听一个 TCP 端口
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
